# example/go-filecoin.conf
[program:go-filecoin]
command=/data/go-filecoin/filecoin/go-filecoin daemon
stderr_logfile=/data/log/go-filecoin.err.log
stdout_logfile=/data/log/go-filecoin.out.log
user=ubuntu